# vendor/common/o2cb.service
[Unit]
Description=Load o2cb Modules
Requires=network.service
After=network.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/o2cb.init enable
ExecStop=/sbin/o2cb.init disable
ExecRestart=/sbin/o2cb.init restart

[Install]
WantedBy=multi-user.target
