# debian/ocfs2-tools.files
sbin/fsck.ocfs2
sbin/mkfs.ocfs2
sbin/mounted.ocfs2
sbin/tunefs.ocfs2
sbin/debugfs.ocfs2
sbin/o2cb_ctl
sbin/mount.ocfs2
sbin/ocfs2cdsl
sbin/ocfs2_hb_ctl
etc/init.d/o2cb
etc/sysconfig/o2cb
